# scripts/gemini-chat-with-update.service
[Unit]
Description=Gemini Chat Application (Auto-Update)
Documentation=https://github.com/Koki-dec/PiChat
After=graphical.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=NODE_ENV=production
WorkingDirectory=/home/pi/ai-chat
ExecStart=/home/pi/ai-chat/scripts/auto-update-start.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# セキュリティ設定
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=graphical.target
